# Library/LaunchAgents/org.partners.passwordAgeCheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>org.partners.passwordAgeCheck</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Library/Application Support/JAMF/Partners/Scripts/ageCheck.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>47400</integer>
    <key>RunAtLoad</key>
    <true/>
  </dict>
</plist>
